# root/etc/systemd/system/keynav.service
[Unit]
Description=Run keynav
Requires=graphical.target

[Service]
ExecStart=/usr/bin/keynav
Environment=DISPLAY=:0
User=leo
Restart=always

[Install]
WantedBy=multi-user.target
